# sv/pgadmin/run
#!/bin/bash

source /etc/envvars

cd /usr/share/pgadmin4/web
exec 2>&1
exec python3 pgAdmin4.py
